# services/base_qos.service
[Unit]
Description=MRDT Base QoS (HTB) setup
After=base_force_vlans.service frr.service
Wants=base_force_vlans.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/pi/Documents/MRDT_Babel_Router/generated_configs/base_qos.sh

[Install]
WantedBy=multi-user.target
