# services/mrdt_dashboard.service
[Unit]
Description=MRDT Babel Router Dashboard
After=network-online.target frr.service
Wants=network-online.target

[Service]
User=root
WorkingDirectory=/home/pi/Documents/MRDT_Babel_Router
Environment=MRDT_DASH_HOST=0.0.0.0
ExecStart=/usr/bin/python3 /home/pi/Documents/MRDT_Babel_Router/dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
